# curator4-1596102745.conf
[program:curator4-1596102745]
command=python3 /home/curator4/curator4-1596102745/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-1596102745/test.err.log
stdout_logfile=/home/curator4/curator4-1596102745/test.out.log